# supervisord.conf
[supervisord]
nodaemon=true
user=root

[program:logpipe]
environment=LC_ALL="C.UTF-8", LANG="C.UTF-8"
command=remote_syslog2 /logs/server /logs/worker syslog://papertail.com
